# headless-docker_0.9.2_arm64/DEBIAN/postinst
#!/bin/bash
systemctl disable mediacenter

